# com.ashu.whisper-hotkey.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ashu.whisper-hotkey</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/osascript</string>
        <string>-e</string>
        <string>tell application "Terminal"
    activate
    do script "python3 /Users/ashu/Documents/Projects/whisper-hotkey/main.py >> /tmp/whisper-hotkey.log 2>&1"
    set miniaturized of window 1 to true
    delay 1
    set miniaturized of window 1 to true
end tell</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/whisper-hotkey-launch.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/whisper-hotkey-launch.log</string>
</dict>
</plist>
